# script/99-mesh-network.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlx*", RUN+="/bin/bash /home/vx/code/vxpollbook/script/setup_basic_mesh.sh"
