# script/99-mesh-network.rules
ACTION=="move", SUBSYSTEM=="net", ENV{ID_MODEL}=="802.11_n_WLAN", RUN+="/bin/bash /vx/scripts/setup_basic_mesh.sh"
